# server/systemd/webdashboard.service
[Unit]
Description=WebDashboard service
After=network.target postgresql.service

[Service]
Type=simple
User=webdashboard
WorkingDirectory=/var/lib/webdashboard/app/server
EnvironmentFile=-/etc/webdashboard/webdashboard.env
Environment=NODE_ENV=production
Environment=PORT=49615
ExecStart=/usr/bin/node /var/lib/webdashboard/app/server/src/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
